# crashplan/files/etc/service/openbox/run
#!/bin/bash
rm -f /tmp/.X1-lock

exec env DISPLAY=:1 HOME=/root \
     /usr/bin/openbox-session