# crashplan/files/etc/service/novnc/run
#!/bin/bash

. /opt/default-values.sh

/opt/novnc/utils/launch.sh --listen ${WEB_PORT} --vnc localhost:${VNC_PORT}